# playbooks/roles/linux-mirror/linux-mirror-systemd/mcgrof-next-mirror.service
[Unit]
Description=Linux mcgrof's mirror [mcgrof-next.git]
Documentation=man:git(1)
ConditionPathExists=/mirror/mcgrof-next.git

[Service]
Type=oneshot
ExecStartPre=/usr/bin/git -C /mirror/mcgrof-next.git remote update --prune
ExecStart=/usr/bin/git -C /mirror/mcgrof-next.git fetch --tags --prune
ExecStartPost=/usr/bin/git -C /mirror/mcgrof-next.git fetch origin +refs/heads/*:refs/heads/*

[Install]
WantedBy=multi-user.target
